# files/common/etc/crontabs/root
* * * * */10 /usr/sbin/p2ptbl gossip /tmp/p2ptbl/ccfg_wifi 500 br-mesh
* * * * *    /usr/sbin/p2ptbl gossip /tmp/p2ptbl/gwip      500 br-mesh
* * * * *    /usr/sbin/p2ptbl gossip /tmp/p2ptbl/splash    500 br-mesh
* * * * *    /usr/sbin/p2ptbl gossip /tmp/p2ptbl/update    500 br-mesh
* * * * * /usr/sbin/fsm watch inetable
* * * * * /usr/sbin/fsm watch update
* * * * * /usr/sbin/splash_sync
