# files/common/etc/crontabs/root
* * * * */10 /usr/sbin/p2ptbl gossip /tmp/p2ptbl/ccfg_wifi 500 br-mesh
*/1 * * * *    /usr/sbin/p2ptbl gossip /tmp/p2ptbl/gwip      500 br-mesh
*/1 * * * *    /usr/sbin/p2ptbl gossip /tmp/p2ptbl/splash    500 br-mesh
*/5 * * * *    /usr/sbin/p2ptbl gossip /tmp/p2ptbl/update    500 br-mesh
*/2 * * * * /usr/sbin/fsm watch inetable
*/5 * * * * /usr/sbin/fsm watch update
*/1 * * * * /usr/sbin/splash_sync

